# supervisord.conf
[supervisord]
nodaemon=true

[program:hatebu]
command=python /opt/kakuyomu_rss_gen/bin/web.py
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
autostart=true
autorestart=true

